# blog.ini
[program:blog]
directory=/home/dherrman/blog
command=node server.js
autostart=true
autorestart=true
environment=NODE_ENV=production,HOSTNAME=0.0.0.0,PORT=3000